# templates/overwatch-resume.service.erb
[Unit]
Description=Overwatch boot-time mode resume (auto-start VM or ollama)
After=multi-user.target libvirtd.service
Wants=libvirtd.service
# Run unconditionally — the resume helper inspects /proc/cmdline itself and
# either starts overwatch@<vm>.service (vm-mode) or ollama.service (host-mode)
# or does nothing (cmdline has no overwatch.mode marker).
ConditionPathExists=/usr/local/bin/overwatch-mode

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/overwatch-mode resume
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
