# templates/overwatch@.service.erb
[Unit]
Description=Overwatch VM lifecycle (instance %i)
After=libvirtd.service
Requires=libvirtd.service

[Service]
Type=simple

# Mutex: refuse to start if another overwatch instance owns /run/overwatch/active-vm.
# acquire-mutex writes %i atomically with O_EXCL; fails non-zero on conflict so
# systemd marks the unit failed and never reaches ExecStart.
ExecStartPre=/usr/local/bin/overwatch acquire-mutex %i

ExecStart=/usr/local/bin/overwatch start %i

# Always release the mutex on stop, regardless of how the unit exited.
ExecStopPost=/usr/local/bin/overwatch release-mutex %i

TimeoutStartSec=infinity
TimeoutStopSec=120
Restart=no

[Install]
WantedBy=multi-user.target
