# templates/overwatch@.service.erb
[Unit]
Description=Overwatch VM lifecycle (instance %i)
After=libvirtd.service
Requires=libvirtd.service

[Service]
Type=simple
<% if @pass4_enabled -%>

# Mode gate: refuse to start if the host isn't booted in vm-mode. If we're
# in host-mode (or any other mode), overwatch-mode stashes %i as the pending
# boot intent, sets the persistent GRUB saved_entry to overwatch-vm-mode (via
# grub-set-default), schedules a reboot, and exits non-zero. systemd fails
# the unit but the reboot is in flight. After the reboot, overwatch-resume.service
# auto-starts overwatch@%i in vm-mode.
ExecStartPre=/usr/local/bin/overwatch-mode require vm %i
<% end -%>

# Mutex: refuse to start if another overwatch instance owns /run/overwatch/active-vm.
# acquire-mutex writes %i atomically with O_EXCL; fails non-zero on conflict so
# systemd marks the unit failed and never reaches ExecStart.
ExecStartPre=/usr/local/bin/overwatch acquire-mutex %i

ExecStart=/usr/local/bin/overwatch start %i

# Always release the mutex on stop, regardless of how the unit exited.
ExecStopPost=/usr/local/bin/overwatch release-mutex %i

TimeoutStartSec=infinity
TimeoutStopSec=120
Restart=no

[Install]
WantedBy=multi-user.target
