# scripts/overwatch.service
[Unit]
Description=Overwatch VM lifecycle (GPU passthrough)
After=libvirtd.service
Requires=libvirtd.service

[Service]
Type=simple
ExecStart=/usr/local/bin/overwatch start
TimeoutStartSec=infinity
TimeoutStopSec=120
KillMode=mixed
Restart=no

[Install]
WantedBy=multi-user.target
